# debian/postinst
#!/bin/sh
set -e

# Create caddy user/group if they don't exist
if ! getent group caddy >/dev/null; then
    groupadd --system caddy
fi

if ! getent passwd caddy >/dev/null; then
    useradd --system \
        --gid caddy \
        --create-home \
        --home-dir /var/lib/caddy \
        --shell /usr/sbin/nologin \
        --comment "Caddy web server" \
        caddy
fi

# Create directories
mkdir -p /etc/caddy
mkdir -p /var/lib/caddy/.local/share/caddy
chown -R caddy:caddy /var/lib/caddy

# Enable and start service
if [ -d /run/systemd/system ]; then
    systemctl daemon-reload
    systemctl enable caddy.service
fi

echo "tengu-caddy installed. Edit /etc/caddy/Caddyfile and run: systemctl start caddy"
